# container-files/etc/profile.d/opt-path.sh
#
# Add /opt/bin to the PATH
#
export PATH=/opt/bin:$PATH
