# container-files-behat/etc/services.d/vnc-server/run
#!/usr/bin/with-contenv /usr/bin/execlineb
foreground { s6-svwait -U /var/run/s6/services/xvfb }
/usr/local/bin/x11vnc -forever -shared -rfbport 5900 -display $DISPLAY -passwd $VNC_PASSWORD
