# container-files-behat/etc/services.d/selenium/run
#!/usr/bin/with-contenv /usr/bin/execlineb
foreground { s6-svwait -U /var/run/s6/services/xvfb }
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
/usr/bin/java -jar /usr/bin/selenium-server-standalone.jar
