# container-files/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv /usr/bin/execlineb
s6-notifyoncheck -s 500 -c "test -S /var/run/php-fpm.sock" /usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.conf -R
