# container-files/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv /usr/bin/execlineb
redirfd -w 2 /dev/null
s6-notifyoncheck /usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.conf -R
