# container-files-behat/etc/services.d/xvfb/run
#!/usr/bin/with-contenv /usr/bin/execlineb
importas display DISPLAY
importas dimension SCREEN_DIMENSION
s6-notifyoncheck /usr/bin/Xvfb $display -ac -screen 0 $dimension
